# client-iso/rootfs-overlay/etc/systemd/system/vpnnet-agent.service
[Unit]
Description=VPN.net Client Agent
Documentation=https://vpn.net/docs/agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/vpnnet-agent.sh
Restart=on-failure
RestartSec=30
StartLimitBurst=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
